# oci-fleet/verifier.service
[Unit]
Description=CafresoAI HQ session verifier (Caddy forward_auth backend)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/cafresoai/oci-fleet
# Shared env file holds HQ_SESSION_SECRET (same hex set on the cafresoai_keys
# canister) and is also sourced by the fleet-api service.
EnvironmentFile=/etc/cafresoai/hq.env
ExecStart=/usr/bin/python3 /opt/cafresoai/oci-fleet/verifier.py
Restart=always
RestartSec=2
# Localhost-only; never expose :9090 publicly.

[Install]
WantedBy=multi-user.target
